# overlay/etc/supervisor.d/avahi.ini
/etc/init.d/avahi-daemon start

[program:avahi]
priority=30
autostart=true
autorestart=true
directory=/
user=%(ENV_USER)s
command=avahi-daemon start
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
stopsignal=INT
stdout_logfile=/var/log/supervisor/avahi.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/avahi.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
